# mediapi.service
[Unit]
Description=MediAPI - Music Player Service
After=network-online.target bluetooth.target pulseaudio.service
Wants=network-online.target
Requires=bluetooth.target

[Service]
# Run as pi user (adjust if needed)
User=pi
Group=pi

# Working directory
WorkingDirectory=/home/pi/Code/mediapi

# Start command - disable auto_connect_bt for systemd to avoid startup delays
ExecStart=/home/pi/Code/mediapi/.venv/bin/python -c "from player import MP3Player; app = MP3Player(use_hardware=True, auto_connect_bt=False); app.run()"

# Restart policy
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mediapi

# Set proper environment for Bluetooth/PulseAudio
Environment="DBUS_SYSTEM_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket"
Environment="PULSE_DBUS_SERVER=unix:path=/run/pulse/dbus"

# Allow access to hardware (GPIO, SPI, etc.)
SupplementaryGroups=gpio spi i2c dialout

# Device access
DevicePolicy=auto

# Type
Type=simple

# Timeouts
TimeoutStartSec=60
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
